# deploy/systemd/stayconnect-edged.service
# edged — Hotel Admin API (Edge API) on the appliance.
# Binds the management loopback only; Caddy exposes it on the management IP.
[Unit]
Description=StayConnect Edge API (Hotel Admin backend)
# Ordering only (best-effort): start after the network, the DB container's
# supervisor, and scd. These are NOT hard requirements — edged tolerates scd
# being down (it dials scd's socket lazily, per request) and retries the site
# DB with backoff at startup instead of exiting, so it must never be stopped or
# blocked just because a dependency is momentarily unavailable.
After=network-online.target docker.service stayconnect-scd.service
Wants=network-online.target docker.service stayconnect-scd.service
# NOTE: intentionally NO Requires= and NO PartOf= scd. A scd restart (runtime
# enrollment, assignment adoption, redeploy) must NOT tear edged down — the
# Hotel Admin API stays continuously available and simply reconnects to scd's
# socket on the next request. This is what fixes the "edged goes inactive when
# scd restarts and never recovers" defect.
# StartLimitIntervalSec=0 (a [Unit] directive) disables the start-rate limiter
# so edged may restart as often as needed and can never be permanently wedged.
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/opt/stayconnect/bin/edged serve
EnvironmentFile=-/etc/stayconnect/edged.env
User=stayconnect
Group=stayconnect
# Always recover, with a bounded 2s backoff. Combined with the in-process DB
# retry and StartLimitIntervalSec=0 above, edged self-heals from any crash,
# dependency delay or socket outage with zero manual intervention.
Restart=always
RestartSec=2
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/stayconnect
ReadOnlyPaths=/etc/stayconnect
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
