# deploy/systemd/stayconnect-scd.service
[Unit]
Description=StayConnect Session Controller Daemon (scd)
# docker.service: the site DB (stayconnect-pg) is a container, so order after
# Docker at cold boot. scd itself uses a lazy pgx pool (no startup Ping), so a
# slow DB does not crash it — this is ordering, not a hard requirement.
After=network-online.target docker.service postgresql.service nftables.service
# portald + acctd now Wants= scd (soft) rather than Require=, so they are not
# stopped when scd stops; Wants= here still pulls them back up when scd starts,
# so a scd stop/start (re-enrollment, assignment-driven restart, redeploy)
# leaves the whole guest plane converged with no manual intervention.
Wants=network-online.target docker.service stayconnect-portald.service stayconnect-acctd.service
Requires=nftables.service
# Never permanently give up (a slow DB/nftables at boot must not wedge the core
# daemon). StartLimitIntervalSec is a [Unit] directive.
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=/etc/stayconnect/scd.env
ExecStart=/opt/stayconnect/bin/scd
# Always recover. scd normally re-execs in place (same PID), so this only fires
# on an actual exit/crash.
Restart=always
RestartSec=2
# needs CAP_NET_ADMIN to run `nft`; running as root for Phase 1 simplicity
User=root
Group=root
RuntimeDirectory=stayconnect
RuntimeDirectoryMode=0755
ProtectSystem=full
# ProtectSystem=full makes /etc read-only, but three stores must persist there:
#   identity/   — the appliance's enrollment identity + signing key
#   license/    — the signed license envelope
#   assignment/ — the SIGNED ASSIGNMENT document: the appliance-local source of
#                 truth for tenant/site. Without this path scd cannot adopt an
#                 assignment ("read-only file system") and stays unassigned.
#   certs/      — the mTLS client certificate + CA bundle
ReadWritePaths=/etc/stayconnect/identity /etc/stayconnect/license /etc/stayconnect/assignment /etc/stayconnect/certs
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
