# deploy/systemd/stayconnect-web-admin.service
[Unit]
Description=StayConnect Admin UI (Next.js 14, production)
After=network-online.target stayconnect-ctrlapi.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/stayconnect/web-admin
Environment=NODE_ENV=production
EnvironmentFile=-/etc/stayconnect/web-admin.env
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=3
User=root
Group=root

[Install]
WantedBy=multi-user.target
