# deploy/systemd/stayconnect-portald.service
[Unit]
Description=StayConnect Captive Portal Daemon (portald)
After=network-online.target stayconnect-scd.service
Wants=network-online.target
Requires=stayconnect-scd.service

[Service]
Type=simple
EnvironmentFile=/etc/stayconnect/portald.env
ExecStart=/opt/stayconnect/bin/portald
Restart=on-failure
RestartSec=2
User=stayconnect
Group=stayconnect
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ProtectSystem=strict
ReadWritePaths=/var/log/stayconnect
ReadOnlyPaths=/etc/stayconnect
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
