# deploy/systemd/stayconnect-netd.service
# netd — privileged network-configuration daemon (Phase 19).
# Root-owned; owns netplan/Kea/nftables/unbound for guest networks. Listens
# only on a protected unix socket. NEVER on TCP.
[Unit]
Description=StayConnect Network Daemon (privileged, guest networks + DHCP)
After=network-online.target docker.service postgresql.service kea-dhcp4-server.service
Wants=network-online.target docker.service
# Never permanently give up: netd owns the guest networks/DHCP/nftables, so it
# must always recover automatically after a crash or a dependency delay.
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/opt/stayconnect/bin/netd
EnvironmentFile=-/etc/stayconnect/netd.env
User=root
Group=root
Restart=always
RestartSec=2
RuntimeDirectory=stayconnect
RuntimeDirectoryMode=0755
# Shared /run/stayconnect also holds scd.sock; do NOT let a restart of any one
# service wipe the others' sockets.
RuntimeDirectoryPreserve=yes
# netd runs ip/nft/netplan/unbound-control and writes /etc/netplan + the
# generated bundle dir; it needs broad filesystem + net admin, so ProtectSystem
# is deliberately relaxed here (this is the one privileged surface).
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
