# deploy/systemd/stayconnect-scd.service.d/hotel-admin-cert.conf
# scd (root) drives the Hotel Admin TLS certificate manager on behalf of edged
# (which is sandboxed with NoNewPrivileges). ProtectSystem=full makes /etc
# read-only; grant write to ONLY the Hotel Admin cert directory so the manager
# can atomically swap the leaf/chain and write its status file there.
[Service]
ReadWritePaths=/etc/caddy/hotel-admin
